# run/builder/runit/pull-buildpack-runner/run
#!/bin/sh

sv start docker || exit 1

docker pull centurylink/buildpack-runner

while true; do sleep 2; done
