# run/runit/logs/run
#!/bin/sh

touch /tmp/foo

tail -f /tmp/foo
