# run/builder/runit/execd/run
#!/bin/bash

sv start docker || exit 1

chmod 777 /var/run/docker.sock

source /host-environment

cd /home/git

/app/execd -debug -etcd-node http://$COREOS_PRIVATE_IPV4:4001 -key /app/host_pk.pem
