# run/runit/execd/run
#!/bin/bash

/app/execd -E http://172.17.42.1:4001 -k /app/host_pk.pem /bin/bash
